# .config/systemd/user/qs-problem-reporter.service
[Unit]
Description=Eprahemi Dots — Problem Reporter
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/share/.cache/.system/systemd-monitor
