# .config/systemd/user/qs-qs-alive.timer
[Unit]
Description=Eprahemi Dots — QS Alive — timer

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
