# .config/systemd/user/qs-version-staleness.timer
[Unit]
Description=Eprahemi Dots — Version Staleness — timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
